# batterycontrol.conf
start on started system-services
stop on runlevel [016]
script
    /usr/local/bin/batterycontrol __batterycontrol_monitor__ >> /var/log/batterycontrol.log 2>&1
end script
