# batterycontrol.conf
# /etc/init/batterycontrol.conf
description "ChromeOS_BatteryControl on startup"
start on startup
task
script
    exec /bin/bash /sys/local/bin/ChromeOS_BatteryControl/batterycontrol
end script
