# batterycontrol.conf
start on started system-services
stop on runlevel [016]

respawn

exec /usr/local/bin/.batterycontrol_conf.sh
